# rootfs/etc/systemd/system/dnsmasq.service.d/override.conf
[Service]
RestartSec=3
StartLimitBurst=0
